# overlay/common/etc/systemd/system/beaconbox-net.service
[Unit]
Description=BeaconBox AP (hostapd + dnsmasq)
After=multi-user.target
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/beaconbox-net
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
